# supervisord.conf
[program:tg-downloader]
directory=/opt/tg-downloader
command=/opt/tg-downloader/venv/bin/python /opt/tg-downloader/bot-enhanced.py
user=tgbot
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/tg-downloader/logs/supervisor.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=PATH="/opt/tg-downloader/venv/bin",HOME="/opt/tg-downloader"

[supervisord]
logfile=/opt/tg-downloader/logs/supervisord.log
pidfile=/var/run/supervisord.pid
nodaemon=false

[supervisorctl]